# scripts/macos/com.novacove.imbued.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.novacove.imbued</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/imbued</string>
        <string>--daemon</string>
        <string>--socket</string>
        <string>~/.imbued/imbued.sock</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>~/.imbued/logs/imbued.log</string>
    <key>StandardOutPath</key>
    <string>~/.imbued/logs/imbued.log</string>
</dict>
</plist>
